# mcp-server/firefox-mcp.service
[Unit]
Description=Firefox MCP Server
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/node /var/home/houge/.openclaw/firefox-mcp/mcp-server/ws-server-v2.js
Restart=always
RestartSec=5

[Install]
WantedBy=default.target
